# system_files/usr/lib/systemd/system/rechunker-group-fix.service
[Unit]
Description=Fix groups
Wants=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=rechunker-group-fix /etc/group
ExecStart=rechunker-group-fix /etc/gshadow

[Install]
WantedBy=default.target multi-user.target
